# rc.local
sleep 5
screen -dmS update-ip /etc/update-ip
exit 0
